# SOURCES/perfmon.service
[Unit]
Description=Performance Monitor Collector
After=network.target

[Service]
Type=simple
ExecStart=/opt/perfmon/bin/perfmon-collector.sh
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
